# watchdog/wireproxy-watchdog.timer
[Unit]
Description=Check wireproxy-wgcf health every 30 seconds

[Timer]
OnBootSec=2min
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
